# scripts/systemd/timeweaver-server.service
[Unit]
Description=TimeWeaver FastAPI Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__TIMEWEAVER_USER__
WorkingDirectory=__TIMEWEAVER_ROOT__/server
Environment=PYTHONUNBUFFERED=1
ExecStart=__TIMEWEAVER_ROOT__/server/.venv/bin/python -m uvicorn app:app --host 0.0.0.0 --port 8000 --workers 1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
